# targets/fs.src/etc_ro_fs/inittab
::sysinit:/etc_ro_fs/rc.d/rc.sysinit
#::askfirst:-/bin/ash
::shutdown:/etc_ro_fs/rc.d/rc.shutdown
::restart:/sbin/init
