# package/beaglebone/postinst
update-rc.d xap defaults 50
/etc/init.d/xap start
